# _config/etc/supervisor/conf.d/expirybot.conf
# copy me into /etc/supervisor/conf.d/ and do `service supervisor reload`
# then start with `supervisorctl start expirybot`

[program:expirybot]
user=expirybot-web
directory=/home/expirybot-web/app
command=/home/expirybot-web/app/script/run_production
redirect_stderr=true
autorestart=true
startsecs=10
# http://docs.gunicorn.org/en/stable/signals.html
stopsignal=QUIT
